# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

# 启动定时任务
exec \
	cron -f -L /dev/null