# root/etc/services.d/subfinder/run
#!/usr/bin/with-contenv bash

# 启动subfinder定时查找字幕

echo "*/60       *       *       *       *       exec \
	s6-setuidgid abc \
	subfinder /media \
	-c /config/subfinder.json" > /var/spool/cron/crontabs/root
