# root/etc/services.d/subfinder/run
#!/usr/bin/with-contenv bash

# 设定subfinder定时查找字幕
if [ `grep  -c subfinder /var/spool/cron/crontabs/root` -eq 0 ];then
echo "*/60       *       *       *       *       exec \
	s6-setuidgid abc \
	subfinder /media \
	-c /config/subfinder.json" > /var/spool/cron/crontabs/root
fi
