# root/etc/services.d/subfinder/run
#!/usr/bin/with-contenv bash

sleep 360d
