# root/etc/services.d/subfinder/run
#!/usr/bin/with-contenv bash

# 启动subfinder查找字幕
exec \
	s6-setuidgid abc \
	subfinder /libraries \
	-c /config/subfinder.json
